# config/ubuntu/upstart/mongod.conf
description "Mongo Database"
author "10Gen"

start on startup
stop on shutdown

exec /usr/local/bin/mongod --logpath /var/log/mongod.log
